# crontab
0 6 * * * bash /root/WorkSpace/HealthPunch/HealthPunch.sh
