# crontab
0 6 * * * bash ~/WorkSpace/HealthPunch/HealthPunch.sh
